# twitter-bot.conf
[program:twitter-bot]
command=/opt/run/twitter-bot/twitter-bot
directory=/opt/run/twitter-bot
user=psi
